# etc/systemd/system/tzupdate.service
[Unit]
Description=Update timezone automatically

[Service]
Type=oneshot
ExecStart=/usr/bin/tzupdate
